# flask-dev/rootfs/etc/s6-overlay/s6-rc.d/celery-beat/finish
#!/usr/bin/with-contenv bashio

bashio::log.green "FINISH: Celery Beat stopped."
