# playit-gg/rootfs/etc/s6-overlay/s6-rc.d/playit/finish
#!/usr/bin/with-contenv bashio

bashio::log.green "FINISH: PlayIt stopped."

exit 0