# playit-gg/rootfs/etc/s6-overlay/s6-rc.d/playit/run
#!/usr/bin/with-contenv bashio

playitd
