# vscode-tunnel/rootfs/etc/s6-overlay/s6-rc.d/code/finish
#!/usr/bin/with-contenv bashio

bashio::log.green "FINISH: VSCode Tunnel stopped."

exit 0