# vscode-tunnel/rootfs/etc/s6-overlay/s6-rc.d/code/run
#!/usr/bin/with-contenv bashio

code tunnel --accept-server-license-terms
