# vscode-tunnel/rootfs/etc/s6-overlay/s6-rc.d/code/run
#!/usr/bin/with-contenv bashio
name=$(bashio::config 'name')
ghLogin=$(bashio::config 'ghLogin')

args=""

if ["$name" = ""]; then
    args="--random-name"
else
    args="--name $name"
fi

tmux new-session -d -s vscode "/code/.code tunnel --accept-server-license-terms $args"

if [ "$ghLogin" != "" ]; then
    tmux send-keys -t vscode "$(printf '\033[B')" # press down arrow to select GitHub login
tmux send-keys -t vscode C-m # press enter to lock in Microsoft or GitHub login

tmux attach-session -t vscode