# vscode-tunnel/rootfs/etc/s6-overlay/s6-rc.d/code/run
#!/usr/bin/with-contenv bashio

name=$(bashio::config 'name')
ghLogin=$(bashio::config 'ghLogin')

bashio::log.green "Name: '$name', GitHub Login: $ghLogin"

checkDir () {
    if [ ! -d "$1" ]; then
        bashio::log.warning "Directory $1 does not exist. Creating..."
        mkdir -p "$1"
    fi
}

checkDir "/config/workspace"
checkDir "/config/.vscode"
checkDir "/config/.vscode-server"
rm -R /root
ln -s /config/root /

args=""

if [ "$ghLogin" = "true" ]; then
    args="github"
else
    args="microsoft"
fi

if [ ! -f "/config/.vscode/token.json" ]; then
    bashio::log.green "No token found, initiating log in with type: $args"
    code tunnel user login --provider $args
fi

if [ "$name" = "" ]; then
    args="--random-name"
else
    args="--name $name"
fi

bashio::log.green "Starting VSCode Tunnel with args: $args"
cd /config/workspace

code tunnel --accept-server-license-terms $args
