# flask/rootfs/etc/s6-overlay/s6-rc.d/celery-beat/run
#!/usr/bin/with-contenv bashio

set -euo pipefail

getConfig() {
  jq --raw-output "$1" /data/options.json
}

cd /config/flask

if [ "$(getConfig '.enableCeleryBeat')" = "true" ]; then
  bashio::log.green "Starting Celery Beat..."
  extraArgs=""
  bashio::log "celeryApp exists: $(bashio::config.exists 'celeryApp')"
  if [ "$(bashio::config.exists 'celeryApp')" = "true" ] && [ "$(getConfig '.celeryApp')" != "" ]; then
    extraArgs="-A $(getConfig '.celeryApp')"
  fi
  bashio::log "Celery Beat args: '${extraArgs}'"

  bashio::log "celeryBeatLogLevel exists: $(bashio::config.exists 'celeryBeatLogLevel')"
  loglevel=$(getConfig '.celeryLogLevel')
  if [ "$(bashio::config.exists 'celeryBeatLogLevel')" = "true" ]; then
    loglevel=$(getConfig '.celeryBeatLogLevel')
  fi
  bashio::log "Celery Beat log level set to '${loglevel}'"

  celery beat --loglevel=${loglevel} ${extraArgs} 2>&1 | sed 's/^/[Celery Beat] /'
else
  bashio::log.info "Celery Beat disabled"
  tail -f /dev/null
fi
