# flask/rootfs/etc/s6-overlay/s6-rc.d/celery-beat/run
#!/usr/bin/with-contenv bashio

set -euo pipefail

getConfig() {
  jq --raw-output "$1" /data/options.json
}

tmpfile="/tmp/celery_beat.log"

cd /config/flask

if [ "$(getConfig '.enableCeleryBeat')" = "true" ]; then
  bashio::log.green "Starting Celery Beat..."
  extraArgs=""
  bashio::log "celeryApp exists: $(bashio::config.exists 'celeryApp')"
  if bashio::config.has_value 'celeryApp'; then
    extraArgs="-A $(getConfig '.celeryApp')"
  fi
  bashio::log "Celery Beat args: '${extraArgs}'"

  bashio::log "celeryBeatLogLevel exists: $(bashio::config.exists 'celeryBeatLogLevel')"
  loglevel=$(getConfig '.celeryLogLevel')
  if bashio::config.true 'celeryBeatLogLevel'; then
    loglevel=$(getConfig '.celeryBeatLogLevel')
  fi
  bashio::log "Celery Beat log level set to '${loglevel}'"

  celery beat --loglevel=${loglevel} ${extraArgs} > ${tmpfile} 2>&1
  tail -n 0 -F "$tmpfile" | while read -r line; do
      echo "[Celery Beat] $line"
  done
else
  bashio::log.info "Celery Beat disabled"
  tail -f /dev/null
fi
