# flask/rootfs/etc/s6-overlay/s6-rc.d/celery-beat/run
#!/usr/bin/with-contenv bashio

set -euo pipefail

getConfig() {
  jq --raw-output "$1" /data/options.json
}

cd /config/flask

if [ "$(getConfig '.enableCeleryBeat')" = "true" ]; then
  bashio::log.green "Starting Celery Beat..."
  extraArgs=""
  if [ $(bashio::config.exists 'celeryApp') = "true" ] && [ "$(getConfig '.celeryApp')" != "" ]; then
    extraArgs="-A $(getConfig '.celeryApp')"
  fi
  loglevel=$(getConfig '.celeryLogLevel')
  if [ $(bashio::config.exists 'celeryBeatLogLevel') = "true" ]; then
    loglevel=$(getConfig '.celeryBeatLogLevel')
  fi
  celery beat --loglevel=${loglevel} ${extraArgs} 2>&1 | sed 's/^/[Celery Beat] /'
else
  bashio::log.info "Celery Beat disabled"
  tail -f /dev/null
fi
