# redis/rootfs/etc/s6-overlay/s6-rc.d/redis-exporter/finish
#!/usr/bin/with-contenv bashio

bashio::log.green "FINISH: Redis Exporter stopped."

exit 0