# redis/rootfs/etc/s6-overlay/s6-rc.d/redis/finish
#!/usr/bin/with-contenv bashio

bashio::log.green "FINISH: Redis stopped."

exit 0