# redis/rootfs/etc/s6-overlay/s6-rc.d/redis/run
#!/usr/bin/with-contenv bashio

bashio::log.green "Starting Redis..."

redis-server --protected-mode no --port 6379
