# flask/rootfs/etc/s6-overlay/s6-rc.d/flask/finish
#!/usr/bin/with-contenv bashio

bashio::log.green "FINISH: Flask stopped."

exit 0