# flask/rootfs/etc/s6-overlay/s6-rc.d/flask/run
#!/usr/bin/with-contenv bashio

set -euo pipefail

getConfig() {
  jq --raw-output "$1" /data/options.json
}

bashio::log.green "Starting app..."

# Start Flask app
python3 /config/flask/main.py