# flask/rootfs/etc/s6-overlay/s6-rc.d/flask/run
#!/usr/bin/with-contenv bashio

set -euo pipefail

getConfig() {
  jq --raw-output "$1" /data/options.json
}

# Activate the virtual environment
bashio::log.green "Activating virtual environment in /config/flask"
source /config/flask/bin/activate

bashio::log.green "Starting app..."
cd /config/flask

python3 ./main.py