# flask/rootfs/etc/s6-overlay/s6-rc.d/celery/run
#!/usr/bin/with-contenv bashio

set -euo pipefail

getConfig() {
  jq --raw-output "$1" /data/options.json
}

cd /config/flask

if [ "$(getConfig '.enableCelery')" = "true" ]; then
  bashio::log.green "Starting Celery..."
  extraArgs=""
  if [ bashio::config.exists 'celeryApp' ] && [ "$(getConfig '.celeryApp')" != "" ]; then
    extraArgs="-A $(getConfig '.celeryApp')"
  fi
  celery worker --loglevel=$(getConfig '.celeryLogLevel') ${extraArgs} 2>&1 | sed 's/^/[Celery] /'
else
  bashio::log.info "Celery disabled"
  tail -f /dev/null
fi
