# mc-all-platform-velocity-beta/rootfs/etc/s6-overlay/s6-rc.d/velocity/finish
#!/usr/bin/with-contenv bashio

bashio::log.green "FINISH: Sending graceful stop command..."
if tmux has-session -t velocity 2>/dev/null; then
  tmux send-keys -t velocity "stop" C-m
  bashio::log.green "FINISH: stop command sent. Waiting for Velocity process to exit..."
  while tmux has-session -t velocity 2>/dev/null; do
    sleep 1
  done
  bashio::log.green "FINISH: Velocity has fully stopped."
else
  bashio::log.green "FINISH: No active tmux session 'velocity' found."
fi
exit 0